# scrtips/supervisor/supervisor.ini
[inet_http_server]
port=127.0.0.1:9001

[program:rdb]
directory=/app/rdb
environment=RAFT_BOOTSTRAP=true
environment=RAFT_JOIN_ADDR=${leader_raft_http_bind}
command=rdb -config ./conf.yaml
stderr_logfile=/app/rdb.log
stdout_logfile=/app/rdb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10